# resources/supervisor/local.conf
[program:telegraf]
command=/usr/bin/telegraf
autostart=false

[program:automaton]
command=automaton lifecycle.yml -d -s /tmp/sock
stopsignal=INT
stopasgroup=true
stopwaitsecs=60